# advance/files/initd_upwifi
#!/bin/sh /etc/rc.common

START=9
STOP=12

boot() {
    sleep 20 &&  wifi up
}
